# scripts/podman/legalogy.container.in
# Legalogy gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into legalogy's ~/.config/containers/systemd/
# {{LEGALOGY_HOME}} is replaced at install time.

[Unit]
Description=Legalogy gateway (rootless Podman)

[Container]
Image=legalogy:local
ContainerName=legalogy
UserNS=keep-id
Volume={{LEGALOGY_HOME}}/.legalogy:/home/node/.legalogy
EnvironmentFile={{LEGALOGY_HOME}}/.legalogy/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
